# xwl@.service.in
[Unit]
Description=Xwayland server at display %i

[Service]
Type=notify
ExecStart=@bindir@/xwl-run --display=%i --no-exit-with-child /bin/sh -c "systemctl --user import-environment DISPLAY XWL_VERSION; @sysconfdir@/xwlrc"
Restart=always

[Install]
WantedBy=default.target
